# api/env/rssapp.conf
[program:rssapp_worker]
directory=/home/ubuntu/projects/rssapp
command=celery worker -A tasks -l info -E
autostart=true
autorestart=true
